# roles/gvl.commandline-utilities/files/jupyterhub.conf
# Upstart configuration for JupyterHub
description "jupyterhub"

start on (filesystem)
stop on runlevel [016]

exec /usr/local/bin/jupyterhub --no-ssl -f /etc/jupyterhub/jupyterhub_config.py